# config/systemd/user/tray.target
[Unit]
Description=Window manager tray
Requires=polybar.service
After=polybar.service
PartOf=polybar.service

[Install]
WantedBy=polybar.service
