# config/systemd/user/wm.target
[Unit]
Description=Window manager target
Requires=graphical-session.target
